# files/files/etc/uci-defaults/12_onion_version
#!/bin/sh

uci -q batch <<-EOF > /dev/null
	set onion.@onion[0].version='0.0.7'
	set onion.@onion[0].build='293'
	commit system
EOF

exit 0
