# files/etc/uci-defaults/12_onion_version
#!/bin/sh

uci -q batch <<-EOF > /dev/null
	set onion.@onion[0].version='1.0.0'
	set onion.@onion[0].build='122'
	commit system
EOF

exit 0
