# files/etc/uci-defaults/16_mcu_flash
#!/bin/sh

# on first boot, flash the microcontroller!

sh /usr/bin/mcu/flash.hex 

exit 0